# services/set-user-id/run
#!/usr/bin/env sh

set -e

usermod -u "$PUID" www-data
groupmod -g "$PGID" www-data